# src/DemoWorkerDocker/demoworkerdocker.service
[Unit]
Description=Demo Worker Docker

[Service]
Type=notify
ExecStart=/app/DemoworkerDocker

[Install]
WantedBy=multi-user.target